# apps/data_processor/data_processor.cron
16 * * * * /usr/local/bin/python -m data_processor >> /var/log/data_processor.log 2>&1

